# priv/ansible/roles/redis_server/files/redis-stack.service
[Unit]
Description=Redis stack server
Documentation=https://redis.io/
After=network.target

[Service]
Type=notify
User=root
Group=root
ExecStart=/opt/redis-stack/bin/redis-stack-server /etc/redis/redis.conf
Restart=on-failure
RestartSec=5
WorkingDirectory=/data
UMask=0077
LimitNOFILE=65535
TimeoutStartSec=60
TimeoutStopSec=60

[Install]
WantedBy=multi-user.target

[Unit]
Description=Advanced key-value store
After=network.target
