# priv/ansible/roles/beam_linux_tuning/files/rc.local
defrt=`ip route | grep "^default" | head -1`
ip route change $defrt initcwnd 10
